# scripts/usecase-agent-run.timer
[Unit]
Description=usecase-agent Hourly Run

[Timer]
OnCalendar=*-*-* *:00:00
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
